# infrastructure/deployment/linux/FormXChange.Api.service
[Unit]
Description=FormXChange API Gateway Service
After=network.target postgresql.service

[Service]
Type=notify
ExecStart=/usr/local/formxchange/api/FormXChange.Api
Restart=always
RestartSec=10
User=formxchange
Group=formxchange
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/formxchange

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=formxchange-api

[Install]
WantedBy=multi-user.target




